# services/monitoring.service
[Unit]
Description=Write log every minute

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring.sh

[Install]
WantedBy=multi-user.target